# crontab
# runs every 5 minutes
*/5 * * * * /usr/local/bin/python /usr/src/app/src/main.py >> /var/log/cron.log 2>&1
